# crontab
# m h  dom mon dow   command
0 5  * * *  /home/admin/luigi-cabonare/cron_job.sh
